# reportes/crontab
# m h  dom mon dow   command
0 * * * * /usr/bin/curl http://localhost/nazareno/reportes/sao.php
